# manomonitor.service
[Unit]
Description=ManoMonitor - WiFi Presence Detection System
After=network.target
Documentation=https://github.com/local-gh0st/ManoMonitor

[Service]
Type=simple
User=manomonitor
Group=manomonitor
WorkingDirectory=/opt/manomonitor
Environment="PATH=/opt/manomonitor/venv/bin"
EnvironmentFile=/opt/manomonitor/.env
ExecStart=/opt/manomonitor/venv/bin/python -m manomonitor.main
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=manomonitor

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/manomonitor/data
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Allow network capture capabilities
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
